# usr/lib/systemd/system/update-motd.timer
[Unit]
Description=Update of the system MOTD

[Timer]
OnBootSec=1min
OnUnitActiveSec=10m

[Install]
WantedBy=multi-user.target

