# etc/update-motd.d/10-header
#!/bin/bash

RESTORE='\033[0m'

W="\033[00;37m"
B="\033[01;36m"
R="\033[01;34m"
X="\033[01;37m"
A="\033[01;32m"

R_BG="\e[41m"


####################################### WHOAMI ($ME:10chars)
ME=$(printf "%-10s\n" $(whoami))
####################################### END WHOAMI

####################################### System uptime ($UP:12chars)
uptime=`cat /proc/uptime | cut -f1 -d.`
upDays=$((uptime/60/60/24))
upHours=$((uptime/60/60%24))
upMins=$((uptime/60%60))
upSecs=$((uptime%60))
UP=$(printf "%-3s %-2s:%-2s:%-2s\n" $upDays $upHours $upMins $upSecs)

####################################### Time of day ($TIME:14chars)
HOUR=$(date +"%H")

if [ $HOUR -lt 12  -a $HOUR -ge 0 ]
then   TIME="morning  "
elif [ $HOUR -lt 17 -a $HOUR -ge 12 ]
then   TIME="afternoon"
else   TIME="evening  "
fi
####################################### END Time of day

####################################### SystemLoad ($LOAD1:4chars;$LOAD5:4chars;$LOAD15:4chars)
LOAD1=`cat /proc/loadavg | awk {'print $1'}`
LOAD5=`cat /proc/loadavg | awk {'print $2'}`
LOAD15=`cat /proc/loadavg | awk {'print $3'}`
####################################### END SystemLoad

####################################### SystemUpdate ($UPDATES:18chars)
updates=$(pacman -Qu | wc -l)

if (( $updates > 0 ))
then    UPDATES="$R_BG!UPDATE AVAILABLE!$RESET"
else    UPDATES="                  "
fi
####################################### END SystemUpdate

echo -e " 
$RESTORE
                                                           /\             /\ \n                                                            |\`\\_,--=\"=--,_//\`|
                                                           \ .\"  :'. .':  \". /
                                                          ==)  _ :  '  : _  (==
                                                            |>/O\   _   /O\<|
                                                            | \-\"~\` _ \`~\"-/ |
                                                           >|\`===. \_/ .===\`|<
                                                     .-\"-.   \==='  |  '===/   .-\"-.
.---------------------------------------------------{'. '\`}---\,  .-'-.  ,/---{.'. '}----------------.
 )                                                  \`\"---\"\`     \`~-===-~\`     \`\"---\"\`               (
(                                                                                                    )
 )       \033[01;32m. \033[01;37m            __        _______ _     ____ ___  __  __ _____          \033[0m    Good $TIME   (
(       \033[01;32m/#\ \033[01;37m           \ \      / / ____| |   / ___/ _ \|  \/  | ____| \033[0m               $ME     )
 )     \033[01;32m/###\ \033[01;37m           \ \ /\ / /|  _| | |  | |  | | | | |\/| |  _|   \033[0m  UPTIME: $UP       (
(     \033[01;32m/#####\ \033[01;37m           \ V  V / | |___| |__| |__| |_| | |  | | |___  \033[0m  LOAD: $LOAD1:$LOAD5:$LOAD15        )
 )   \033[01;32m/##.-.##\ \033[01;37m           \_/\_/  |_____|_____\____\___/|_|  |_|_____| \033[0m                             (
(   \033[01;32m/##(   )##\ \033[01;37m\033[01;34m                                                 _     \033[0m                              )
 ) \033[01;32m/#.--   --.#\ \033[01;37m\033[01;34m                                  _ __  ___ _ _| |_ _  _              \033[0m             (
( \033[01;32m/'           '\ \033[01;36m\033[01;34m                                | '  \/ _ \ ' \  _| || |             \033[0m              )
 )\033[01;34m                                                |_|_|_\___/_||_\__|\_, |  \033[0m                        (
(\033[01;34m                                                                    |__/   \033[0m                         )
 )                                                                               $UPDATES (
'----------------------------------------------------------------------------------------------------'
"

